# roles/rr-stream-record/templates/stream_recorder.timer.j2
[Unit]
Description=Run stream recorder every {{ record_overlap }} seconds before the hour

[Timer]
OnCalendar={{ record_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
