# roles/rr-stream-record/files/stream_recorder.service
[Unit]
Description=Stream Recorder Script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/record_stream.sh